# deploy/mtgc.service
[Unit]
Description=MTG Card Collection Web Server
After=network.target

[Service]
Type=simple
User=mtgc
Group=mtgc
WorkingDirectory=/opt/mtgc
Environment=MTGC_HOME=/var/lib/mtgc
ExecStart=/opt/mtgc/.uv/bin/uv run mtg crack-pack-server --port 8081
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
